# packaging/deb/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        if ! getent group planmgrgrp >/dev/null; then
            if getent group 1002 >/dev/null; then
                echo "planmgr: gid 1002 is already in use by another group; cannot create planmgrgrp" >&2
                exit 1
            fi
            groupadd --system --gid 1002 planmgrgrp
        fi
        if ! getent passwd planmgruser >/dev/null; then
            if getent passwd 1002 >/dev/null; then
                echo "planmgr: uid 1002 is already in use by another user; cannot create planmgruser" >&2
                exit 1
            fi
            useradd --system --uid 1002 --gid planmgrgrp --home /var/planmgr --shell /usr/sbin/nologin planmgruser
        fi
        if getent group docker >/dev/null; then
            usermod -aG docker planmgruser
        fi

        install -d -o planmgruser -g planmgrgrp -m 750 /etc/planmgr
        install -d -o planmgruser -g planmgrgrp -m 700 /etc/planmgr/secrets
        install -d -o planmgruser -g planmgrgrp -m 750 /var/planmgr
        install -d -o planmgruser -g planmgrgrp -m 750 /var/planmgr/export
        install -d -o planmgruser -g planmgrgrp -m 750 /var/log/planmgr

        set -a
        . /etc/default/planmgr
        set +a

        if [ -z "$PLANMGR_IMAGE_VERSION" ]; then
            echo "planmgr: PLANMGR_IMAGE_VERSION is not set in /etc/default/planmgr" >&2
            exit 1
        fi

        if [ -z "${PLANMGR_INSTANCE_UUID:-}" ]; then
            if [ ! -s /etc/planmgr/instance_uuid ]; then
                python3 -c 'import uuid; print(uuid.uuid4())' > /etc/planmgr/instance_uuid
                chown planmgruser:planmgrgrp /etc/planmgr/instance_uuid
                chmod 640 /etc/planmgr/instance_uuid
            fi
            PLANMGR_INSTANCE_UUID="$(cat /etc/planmgr/instance_uuid)"
            export PLANMGR_INSTANCE_UUID
        fi

        # TLS material is a secret and is NOT shipped in this package. The
        # operator must provision the following files into /etc/planmgr/secrets
        # out-of-band (e.g. via scp or a secrets manager) before the service
        # can start under the https/mtls protocols:
        #   ca.crt, server.crt, server.key, client.crt, client.key
        # Verify their presence and abort the configure step with an explicit
        # message if any is missing, rather than starting a server that would
        # crash-loop on a missing certificate.
        missing_secret=0
        for f in ca.crt server.crt server.key client.crt client.key; do
            if [ ! -s "/etc/planmgr/secrets/$f" ]; then
                echo "planmgr: required TLS file /etc/planmgr/secrets/$f is missing;" \
                     "provision it out-of-band before starting the service" >&2
                missing_secret=1
            fi
        done
        if [ "$missing_secret" -ne 0 ]; then
            echo "planmgr: TLS material incomplete — rendering config and enabling" \
                 "the service, but NOT starting it. Place the certificates in" \
                 "/etc/planmgr/secrets and run: systemctl start planmgr.service" >&2
        fi

        if ! docker manifest inspect "$PLANMGR_IMAGE_REPO:$PLANMGR_IMAGE_VERSION" >/dev/null; then
            echo "planmgr: image $PLANMGR_IMAGE_REPO:$PLANMGR_IMAGE_VERSION was not found on the registry" >&2
            exit 1
        fi
        docker pull "$PLANMGR_IMAGE_REPO:$PLANMGR_IMAGE_VERSION"

        envsubst < /etc/planmgr/config.json.template > /etc/planmgr/config.json
        chown planmgruser:planmgrgrp /etc/planmgr/config.json
        chmod 640 /etc/planmgr/config.json

        if [ ! -f /etc/planmgr/secrets/db_password ]; then
            # First-run database initialization happens inside the
            # container entrypoint, idempotently, once the service below
            # starts the container; it consumes this generated password
            # itself. This script only seeds the secret, it never runs
            # database initialization directly.
            ( umask 077 && head -c 32 /dev/urandom | base64 > /etc/planmgr/secrets/db_password )
            chown planmgruser:planmgrgrp /etc/planmgr/secrets/db_password
        fi

        systemctl daemon-reload
        systemctl enable planmgr.service
        # Starting the service starts the container; the container
        # entrypoint runs the idempotent database initialization inside
        # the container before the server process starts. This script
        # never invokes database initialization directly. The service is
        # only started when the out-of-band TLS material is present; when
        # it is missing the unit is left enabled but stopped so the
        # operator can provision the secrets and start it explicitly.
        if [ "$missing_secret" -eq 0 ]; then
            systemctl restart planmgr.service
        fi
        ;;
esac

exit 0
