# packaging/deb/prerm
#!/bin/sh
set -e

systemctl stop planmgr.service || true
systemctl disable planmgr.service || true
docker rm -f planmgr || true

exit 0
